# tools/zfs-scrub.service
[Unit]
Description=ZFS scrub checker (all pools)

[Service]
Type=oneshot
ExecStart=/opt/45drives/tools/zfs-scrub-check.sh
